# containers/influxdb/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:influxdb]
user=root
command=/bin/bash -c "./opt/influxdb/versions/0.7.3/influxdb -config=/opt/influxdb/current/config.toml"
