# containers/influxdb/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:influxdb]
command=/bin/bash -c "service influxdb start"
