# containers/shinken/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:shinken]
command=/bin/bash -c "service shinken start"

[program:ssh]
command=/bin/bash -c "service ssh start"
