# containers/rekishi/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/bin/bash -c "service ssh start"

[program:adagios]
command=/bin/bash -c "cd /opt/rekishi/rekishi && python manage.py runserver 0.0.0.0:8000"
